# deploy/systemd/nijenhuis-frontend.service
[Unit]
Description=Nijenhuis Frontend Development Server (Vite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=andre
Group=andre
WorkingDirectory=/home/andre/Desktop/Projects/Nijenhuis/frontend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=development"

# Start Vite dev server
ExecStart=/usr/bin/npm run dev

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/home/andre/Desktop/Projects/Nijenhuis/logs/frontend_server.log
StandardError=append:/home/andre/Desktop/Projects/Nijenhuis/logs/frontend_server.log

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=100%

[Install]
WantedBy=multi-user.target
